# ansible/physical/roles/libvirtd/files/etc_systemd_system/libvirtd.service.d/override.conf
# Execute `systemctl cat libvirtd` to see the complete unit.
[Unit]
Wants=virtlogd.service
